# debian/DEBIAN/postinst
#!/bin/bash

systemctl enable hobot-display.service
exit 0
