# debian/DEBIAN/prerm
#!/bin/bash

systemctl disable hobot-display.service

exit 0
