# debian/DEBIAN/prerm
#!/bin/bash
update-rc.d x3dispinit remove

systemctl disable x3dispinit.service
exit 0

